# roles/airflow/templates/airflow.service.j2
[Unit]
Description=Airflow Webserver

[Service]
ExecStart=airflow webserver -p 8080
User=airflow
Restart=on-failure

[Install]
WantedBy=multi-user.target