# roles/nats/templates/nats.service.j2
[Unit]
Description=NATS messaging server

[Service]
ExecStart=/opt/gnatsd -c /etc/opt/nats/gnatsd.config
User=nats
Restart=on-failure

[Install]
WantedBy=multi-user.target