# roles/redis/templates/redis.service.j2
[Unit]
Description=Redis

[Service]
ExecStart=/opt/redis-server
User=redis
Restart=on-failure

[Install]
WantedBy=multi-user.target