# roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus

[Service]
ExecStart=/opt/prometheus --config.file=prometheus.yml 
User=prometheus
Restart=on-failure

[Install]
WantedBy=multi-user.target